# shell/crontab/test.sh
#!/bin/bash
# 生成的aa 在liutailin用户的home目录

echo `date +%Y-%m-%d\ %H:%M:%S` >> aa
